# book_pipeline/com.textbookreader.bookpipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  textbook-reader 教科書 pipeline 自動化迴圈 launchd agent（standby，每天 08:30 台灣一次 tick）。
  節奏=daily 對齊額度 reset（zlib/MinerU 每日 UTC 0:00 reset；08:30 台灣=UTC 00:30，reset 後）。
  一個 daily tick 把當日額度榨乾：crawl 爬滿 zlib 10/日、ingest 打滿 MinerU、audit/catalog 全做完。
  安裝（需使用者授權，會建立持久背景排程 + RunAtLoad 立即跑一次全鏈 tick）：
    cp book_pipeline/com.textbookreader.bookpipeline.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.textbookreader.bookpipeline.plist
  停用：
    launchctl unload ~/Library/LaunchAgents/com.textbookreader.bookpipeline.plist
  注意：路徑寫死 /Users/chenliangyu；換機需改。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.textbookreader.bookpipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/chenliangyu/project/textbook-reader/book_pipeline/daemon_run.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/chenliangyu/project/textbook-reader/book_pipeline/reports/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chenliangyu/project/textbook-reader/book_pipeline/reports/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/chenliangyu</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
