# book_pipeline/com.textbookreader.bookpipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  textbook-reader 教科書 pipeline 自動化迴圈 launchd agent（standby，反應式控制迴圈）。
  節奏=StartInterval 900s（15min）重拉。daemon_run.sh 設 BOOK_PIPELINE_REACTIVE=1：一個
  controller 進程跑有界 observe→非阻塞派工→reap→harvest→sleep 迴圈，三條件齊備（產物就緒 ∧
  資源可用 ∧ 無人在做）的 transition 立即派 thread worker。MinerU OCR 為 async：一 poll 到
  就緒即收割（延遲塌縮到一 poll cycle，免人工 kick）。crawl 每 invocation 一次、zlib 額度 0
  即 cheap-skip → 天然退化日節奏（不必固定時刻）。controller 排空或達牆鐘(50min)退出，launchd
  下個 900s 重拉；flock 序列化：上個 controller 未結束則本次 fire 被跳過（不並發同 Label）。
  ⚠ devctl TICK_INTERVAL_S 須與此 900 一致（改一邊改另一邊）。要回滾 daily 單次 tick：
  daemon_run.sh 設 BOOK_PIPELINE_REACTIVE=0 並把此處改回 StartCalendarInterval。
  安裝（需使用者授權，會建立持久背景排程 + RunAtLoad 立即跑一次全鏈 tick）：
    cp book_pipeline/com.textbookreader.bookpipeline.plist ~/Library/LaunchAgents/
    launchctl unload ~/Library/LaunchAgents/com.textbookreader.bookpipeline.plist 2>/dev/null
    launchctl load ~/Library/LaunchAgents/com.textbookreader.bookpipeline.plist
  停用：
    launchctl unload ~/Library/LaunchAgents/com.textbookreader.bookpipeline.plist
  注意：路徑寫死 /Users/chenliangyu；換機需改。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.textbookreader.bookpipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/chenliangyu/project/textbook-reader/book_pipeline/daemon_run.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/chenliangyu/project/textbook-reader/book_pipeline/reports/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chenliangyu/project/textbook-reader/book_pipeline/reports/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/chenliangyu</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
